# lib/conf/linux/prey-agent.conf
description  "Prey Anti Theft Agent"
author       "Tomás Pollak <tomas@forkhq.com>"

start on startup
stop on shutdown
kill signal QUIT

respawn
respawn limit 20 5

env TERM=dumb

script
  PREY_USER={{user}}
  # chown $PREY_USER: /var/log/prey.log
  chdir {{root_path}}
  exec sudo -u $PREY_USER {{prey_bin}}
end script
