# lib/conf/linux/prey-agent.conf
description  "Prey Anti Theft Agent"
author       "Tomás Pollak <tomas@forkhq.com>"

setuid {{user}}
start on startup
stop on shutdown
kill signal QUIT

respawn
respawn limit 20 5

env TERM=dumb

# upstart cookbook: http://upstart.ubuntu.com/cookbook/#setuid

script
  # PREY_USER={{user}}
  # chown $PREY_USER: /var/log/prey.log
  chdir {{root_path}}
  # exec sudo -u $PREY_USER {{prey_bin}}
  # exec start-stop-daemon --start -c $PREY_USER --exec {{prey_bin}}
  exec {{prey_bin}}
end script
