# lib/conf/linux/prey-agent.conf
description  "Prey Anti Theft Agent"
author       "Tomás Pollak <tomas@forkhq.com>"

setuid {{user}}
start on startup
stop on shutdown
kill signal QUIT

respawn
respawn limit 20 5

env TERM=dumb

script
  chdir {{root_path}}
  exec {{prey_bin}}
end script
