# lib/conf/linux/upstart.conf
description "Prey Anti-Theft Agent"
author      "Tomás Pollak"

start on startup
stop on shutdown
respawn
kill signal QUIT

env TERM=dumb

script
  PREY_USER={{user}}
  exec sudo -u $PREY_USER {{prey_bin}}
end script
